# scripts-termux/shortcuts/crontab/crontab_disable.sh
#!/data/data/com.termux/files/usr/bin/bash

echo "Disabling/removing crontab entries..."
crontab -r
crontab -l


